# kismet-webui.service
[Unit]
Description=Kismet WebUI Service
After=network.target

[Service]
Type=exec
User=root
Group=root
WorkingDirectory=/opt/kismet-webui
Environment="PATH=/opt/kismet-webui/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
# IMPORTANT: Generate a secure secret with: python -c "import secrets; print(secrets.token_hex(32))"
Environment="SESSION_SECRET=change-this-to-a-random-secret-key"
Environment="DATABASE_URL=sqlite:////opt/kismet-webui/instance/kismet_webui.db"
ExecStart=/opt/kismet-webui/venv/bin/gunicorn --bind 0.0.0.0:2502 --reuse-port --workers 1 main:app
Restart=always
RestartSec=10

# Security settings
PrivateTmp=true
NoNewPrivileges=false
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/kismet-webui /etc/kismet /var/log/kismet /home

[Install]
WantedBy=multi-user.target